# scripts/valtgeist-alerts.service
[Unit]
Description=Valtgeist alert fleet (cascade / toxic-flow / halt -> Telegram)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/valtgeist/podrunner
# Secrets + tuning live here, NOT in the repo (matches the repo's env-only token rule).
EnvironmentFile=/etc/valtgeist/alerts.env
ExecStart=/opt/valtgeist/.venv/bin/python3 -u /opt/valtgeist/podrunner/fleet.py
Restart=always
RestartSec=10
# light hardening
NoNewPrivileges=true
ProtectSystem=full
# resource guards for a 2GB box: a leak degrades to a service restart, never an OOM'd host;
# the ws feed races many endpoints x pools, so give file descriptors real headroom
MemoryMax=400M
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
